# deploy/coraza-waf-mod.service
# Reference template for a manual/hand-rolled install. `deploy/install.sh`
# does NOT copy this file — it generates its own unit (with flags chosen
# interactively, e.g. TLS on/off) directly on the target host. Edit this file
# only if you're deploying without install.sh.
[Unit]
Description=Coraza WAF Mod (WAF + reverse proxy)
After=network.target

[Service]
Type=simple
User=coraza-waf-mod
Group=coraza-waf-mod
WorkingDirectory=/var/lib/coraza-waf-mod
ExecStart=/usr/local/bin/coraza-waf-mod --listen :80 --listen-tls :443 --db /var/lib/coraza-waf-mod/waf.db --certs /var/lib/coraza-waf-mod/certs --retention 30
Restart=on-failure
RestartSec=5s

# Binds :80/:443 without running as root: grant just the one capability
# needed for privileged ports instead of full root.
AmbientCapabilities=CAP_NET_BIND_SERVICE
CapabilityBoundingSet=CAP_NET_BIND_SERVICE

# Sandboxing: read-only system dirs, private /tmp, no privilege escalation.
# /var/lib/coraza-waf-mod (waf.db, certs/, GeoLite2 db) stays writable since
# ProtectSystem=full only locks down /usr, /boot and /etc.
ProtectSystem=full
PrivateTmp=true
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
